# debian/prerm
deb-systemd-helper disable monitoring-agent.service >/dev/null || true